# etc/systemd/camera.service
[Unit]
Description=Start the camera service

[Service]
WorkingDirectory=/usr/local/src
ExecStart=/usr/bin/python3 /usr/local/src/main.py
SyslogIdentifier=SRV_CAMERA
Restart=always
RestartSec=30

[Install]
After=web.service
WantedBy=multi-user.target